# tilemill.conf
[program:tilemill]
user=root
command=/usr/share/tilemill/index.js --server=true
autorestart=true
stopsignal=INT
