# tilemill.conf
[program:tilemill]
user=root
command=/usr/bin/nodejs /usr/share/tilemill/index.js --server=true
autorestart=true
stopsignal=INT

